# devel/etc/catkin/profile.d/30.fetch_tools.bash
. "/home/neeraj/FetchSim/src/fetch_tools/setup.bash"
